# packaging/systemd/upclink-vpn.service.in
[Unit]
Description=Servicio del sistema UPClink VPN
After=dbus.service network.target

[Service]
Type=dbus
BusName=org.upclink.VPN
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/upclink-vpn-service
User=root
Group=root
Restart=on-failure
RestartSec=2
KillMode=control-group
PrivateTmp=true
ProtectHome=true
UMask=0077
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/etc/resolv.conf
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelTunables=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
